# deploy/systemd/claude-sentinel.timer
[Unit]
Description=Run Claude Sentinel every 3 minutes

[Timer]
OnActiveSec=1min
OnUnitInactiveSec=3min
AccuracySec=30sec
Persistent=true
Unit=claude-sentinel.service

[Install]
WantedBy=timers.target
